# autoinstall/files/systemd/kiosk.service
# Purple Computer Kiosk Service
# This systemd service ensures the Purple Computer environment
# starts automatically and restarts if it crashes
# NOTE: Currently not used - we use getty auto-login instead

[Unit]
Description=Purple Computer Kiosk Mode
After=getty@tty1.service
Wants=getty@tty1.service

[Service]
Type=simple
User=purple
Group=purple
WorkingDirectory=/home/purple
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/purple/.Xauthority"

# Start the Purple Computer environment
# This assumes .xinitrc launches Kitty with the REPL
ExecStart=/usr/bin/startx

# Restart on failure
Restart=on-failure
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/purple/.purple /home/purple/.config

[Install]
WantedBy=multi-user.target
